# udev/etc/udev/rules.d/85-monitor-hotplug.rules
ACTION=="change",\
KERNELS=="card0",\
SUBSYSTEMS=="drm",\
SYMLINK+="testyface",\
RUN+="/usr/bin/systemctl start --user hotplug-monitor.service"
