# autostart/.config/autostart/chrome.desktop
[Desktop Entry]
Type=Application
Name=Google Chrome
Exec=/usr/bin/google-chrome
Terminal=false
