# autostart/.config/autostart/monitor-layout.desktop
[Desktop Entry]
Type=Application
Name=Set the Monitor Layout
Exec=/home/sh/.local/bin/hotplug-monitor
Terminal=false
